# hockey_squad_scraper/deploy/hockey_squad_scraper.service
[Unit]
Description=Hockey Squad Scraper
After=network.target

[Service]
User=<SYSTEM_USER>
Group=<SYSTEM_GROUP>

WorkingDirectory=<ABSOLUTE_PATH_TO>/hockey_squad
EnvironmentFile=<ABSOLUTE_PATH_TO>/hockey_squad/.env
ExecStart=<ABSOLUTE_PATH_TO>/hockey_squad/venv/bin/python -m hockey_squad_scraper.runner

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target